# adm/templates/plugins/archive/{{cookiecutter.name}}/crontab
# keep files 5 days
0 3 * * * {% raw %}{{MFDATA_HOME}}{% endraw %}/bin/cronwrap.sh --lock --low "find {% raw %}{{MODULE_RUNTIME_HOME}}{% endraw %}/var/archive/ -type f -mtime +5 -exec rm -f {} \;" >/dev/null 2>&1

# remove empty directories
0 4 * * * {% raw %}{{MFDATA_HOME}}{% endraw %}/bin/cronwrap.sh --lock --low "find {% raw %}{{MODULE_RUNTIME_HOME}}{% endraw %}/var/archive/ -type d -exec rmdir {} \;" >/dev/null 2>&1
